# deploy/systemd/workshop-helmsman.service
[Unit]
Description=Workshop Helmsman — milestone tracker (FastAPI/uvicorn)
After=network.target

[Service]
Type=simple
User=sai
WorkingDirectory=/Users/sai/workshop-helmsman
Environment="BIND_HOST=0.0.0.0"
Environment="BIND_PORT=8001"
ExecStart=/Users/sai/workshop-helmsman/venv/bin/python -m src
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
